# ops/systemd/arbitrage-backup.timer
[Unit]
Description=Run arbitrage PostgreSQL backup every six hours

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
