# ops/systemd/arbitrage-watchdog.timer
[Unit]
Description=Run arbitrage readiness watchdog

[Timer]
OnBootSec=2m
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
